# board/services/crop-web.service
[Unit]
Description=Crop AI Web Dashboard
After=network.target

[Service]
WorkingDirectory=/root/demo
ExecStart=/usr/bin/python3 /root/demo/crop_web.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
